# etc/profile.d/cow.sh
#!/bin/sh

fortune -s | cowsay -f /usr/share/cowsay/cows/default.cow && echo
